# tools/systemd/carnine.socket
[Unit]
Description=CarNiNe Service Socket
PartOf=foo.service

[Socket]
ListenStream=/tmp/CarNiNe.sock

[Install]
WantedBy=sockets.target